# devpi_db/custom/etc/supervisor/conf.d/postgres.conf
[program:postgres]
user=root
command=/opt/start-postgis.sh
autorestart=true
stopsignal=INT
